# momir.service
[Unit]
Description=Physical Momir Magic Card Printer
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/physical_momir
ExecStart=/home/pi/.local/bin/uv run main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
